# rootfs/etc/s6-overlay/s6-rc.d/aria2-prepare/run
/etc/s6-overlay/scripts/aria2-prepare.sh